# service-python/authserver/conf/supervisord/supervisord.conf
[supervisord]

[supervisorctl]

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /opt/authServer/conf/uwsgi/uwsgi.ini

[program:nginx-app]
user=root
command = /usr/sbin/nginx
